# templates/minecraft.service.erb
[Unit]
Description=Minecraft Server

[Service]
WorkingDirectory="<%= @install_dir %>"
User="<%= @user %>"

ExecStart=/usr/bin/screen -DmS survival-mc java -server -Xmx<%= @heap_start %> -Xms <%= @heap_min %> -jar minecraft_server.jar nogui

ExecStop=/usr/bin/screen -p 0 -S mc -X eval 'stuff "say SERVER SHUTTING DOWN. Saving map..."\\015'
ExecStop=/usr/bin/screen -p 0 -S mc -X eval 'stuff "save-all"\\015'
ExecStop=/usr/bin/screen -p 0 -S mc -X eval 'stuff "stop"\\015'
ExecStop=/bin/sleep 2
RestartSec=10s

Restart=always


[Install]
WantedBy=multi-user.target
